# scripts/com.rohanshah.kite-ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rohanshah.kite-ingest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/rohanshah/Desktop/AI Instructions/Trading Developments/05 - Backtest Engine (separate)/scripts/run_kite_ingest.py</string>
        <string>--days</string>
        <string>2</string>
    </array>

    <!-- Run weekdays at 16:30 IST (after market close) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/rohanshah/Desktop/AI Instructions/Trading Developments/05 - Backtest Engine (separate)/results/kite_ingest_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rohanshah/Desktop/AI Instructions/Trading Developments/05 - Backtest Engine (separate)/results/kite_ingest_stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
